# .config/systemd/user/random-wallpaper.timer
[Unit]
Description=Run random wallpaper script every 30 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
